# src/odcommon/etc/init.d/odloadconf
#! /bin/sh
### BEGIN INIT INFO
# Provides:          odconf
# Required-Start:    udev networking
# Required-Stop:
# Should-Start:      glibc
# Default-Start:     S
# Default-Stop:      0 6
# Short-Description: Load opendomo configuration files
# Description:       Load opendomo configuration files
#
### END INIT INFO
### Copyright(c) 2011 OpenDomo Services SL. Licensed under GPL v3 or later


. /lib/init/vars.sh
. /lib/lsb/init-functions

do_start () {
	CSTCONF="/mnt/odconf/sysconf/cstconf.tar"
	DEFCONF="/mnt/odconf/sysconf/defconf.tar"
	LOGFILE="/var/opendomo/log/loadcfg.log"
	MNGFILE="/usr/local/bin/manage_conf.sh"
	touch $LOGFILE

	# Configure time with ntpdate, RaspberryPi don't have hardware clock
	if text -x "/usr/sbin/ntpdate-debian"; then
		log_action_begin_msg "Configuring date with ntp"
		ntpdate-debian >/dev/null 2>/dev/null
		log_action_end_msg $?
	else
		 log_warning_msg "ntp is not installed, system date may be wrong"
	fi

	# Extract custon config or default, if don't exist anyone create default file
	if test -f $CSTCONF.gz; then
		log_action_begin_msg "Updating opendomo configuration files (custom)"
		$MNGFILE load custom &>$LOGFILE
		log_action_end_msg $?

	elif test -f $DEFCONF.gz; then
		log_action_begin_msg "Updating opendomo configuration files (default)"
		$MNGFILE load default &>$LOGFILE
		log_action_end_msg $?

	else
		log_action_begin_msg "Creating default configuration file"
		# Copy default config files
                $MNGFILE copy &>$LOGFILE

		# Creating default configuration file
                cd /etc/opendomo
                tar cfp $DEFCONF * --owner 1000 --group 1000 2>$LOGFILE && gzip -f $DEFCONF 2>$LOGFILE
                log_action_end_msg $?
	fi

	# Fix perms
	chown -R admin:admin /usr/local/opendomo
	chown -R admin:admin /var/opendomo
	chown -R admin:admin /etc/opendomo
	chmod -R 744 /usr/local/opendomo
	chmod -R 744 /var/opendomo
	chmod -R 744 /etc/opendomo

	# Create wrapper scripts
	log_action_begin_msg "Creating opendomo wrapper scripts"
	su -c /usr/local/bin/createwrappers.sh admin 2>$LOGFILE
	log_action_end_msg $?
}

case "$1" in
  start|"")
	do_start
        ;;
  restart|reload|force-reload)
        echo "Error: argument '$1' not supported" >&2
        exit 3
        ;;
   status)
        echo "Error: argument '$1' not supported" >&2
        exit 3
        ;;
 stop)
        # No-op
        ;;
  *)
        echo "Usage: odloadconf [start]" >&2
        exit 3
        ;;
esac

:
