# src/odcommon/etc/init.d/odloadconf
#! /bin/sh
### BEGIN INIT INFO
# Provides:          odconf
# Required-Start:    udev
# Required-Stop:
# Should-Start:      glibc
# Default-Start:     S
# Default-Stop:      0 6
# Short-Description: Load opendomo configuration files
# Description:       Load opendomo configuration files
#
### END INIT INFO
### Copyright(c) 2011 OpenDomo Services SL. Licensed under GPL v3 or later


. /lib/init/vars.sh
. /lib/lsb/init-functions

do_start () {
	DEFCONF=/mnt/odconf/sysconf/defconf.tar.gz
	CSTCONF=/mnt/odconf/sysconf/cstconf.tar.gz
	LOGFILE="/var/opendomo/log/loadcfg.log"
	CGIDIR="/var/www/cgi-bin/"

	# Selecting config
	if test -f $CSTCONF; then
		CONFIG=$CSTCONF
		log_action_begin_msg "Updating opendomo configuration files"
	else
		if test -f $DEFCONF; then
			CONFIG=$DEFCONF
			log_action_begin_msg "Updating opendomo configuration files (default)"
		else
			exit 1
		fi
	fi

	/usr/bin/logger -p notice "Loading configuration file"
	cd /
	tar xfzp $CONFIG

	# Fix perms
	chown root:root etc
	chmod ugo+rx etc
	chown -R admin:admin /etc/opendomo
	chmod -R 744 /etc/opendomo
	chown -R admin:admin /var/opendomo
	chmod -R 744 /var/opendomo
}

case "$1" in
  start|"")
	do_start
	log_action_end_msg $?
        ;;
  restart|reload|force-reload)
        echo "Error: argument '$1' not supported" >&2
        exit 3
        ;;
   status)
        echo "Error: argument '$1' not supported" >&2
        exit 3
        ;;
 stop)
        # No-op
        ;;
  *)
        echo "Usage: odloadconf [start]" >&2
        exit 3
        ;;
esac

:
