# src/odcommon/etc/init.d/odloadconf
#! /bin/sh
### BEGIN INIT INFO
# Provides:          odconf
# Required-Start:    udev
# Required-Stop:
# Should-Start:      glibc
# Default-Start:     S
# Default-Stop:      0 6
# Short-Description: Load opendomo configuration files
# Description:       Load opendomo configuration files
#
### END INIT INFO
### Copyright(c) 2011 OpenDomo Services SL. Licensed under GPL v3 or later


. /lib/init/vars.sh
. /lib/lsb/init-functions

do_start () {
	CSTCONF="/mnt/odconf/sysconf/cstconf.tar"
	DEFCONF="/mnt/odconf/sysconf/defconf.tar"
	LOGFILE="/var/opendomo/log/loadcfg.log"
	MNGFILE="/usr/local/bin/manage_conf.sh"
	touch $LOGFILE

	# Extract custon config or default, if don't exist anyone create default file
	if test -f $CSTCONF.gz; then
		log_action_begin_msg "Updating opendomo configuration files (custom)"
		$MNGFILE load custom &>$LOGFILE
		log_action_end_msg $?

	elif test -f $DEFCONF.gz; then
		log_action_begin_msg "Updating opendomo configuration files (default)"
		$MNGFILE load default &>$LOGFILE
		log_action_end_msg $?

	else
		log_action_begin_msg "Creating default configuration file"
		# Copy default config files
                $MNGFILE copy &>$LOGFILE

		# Creating default configuration file
                cd /etc/opendomo
                tar cfp $DEFCONF * --owner 1000 --group 1000 2>$LOGFILE && gzip -f $DEFCONF 2>$LOGFILE
                log_action_end_msg $?
	fi

	# Create wrapper scripts
	log_action_begin_msg "Creating opendomo wrapper scripts"
	su -c /usr/local/bin/createwrappers.sh admin 2>$LOGFILE
	log_action_end_msg $?
}

case "$1" in
  start|"")
	do_start
        ;;
  restart|reload|force-reload)
        echo "Error: argument '$1' not supported" >&2
        exit 3
        ;;
   status)
        echo "Error: argument '$1' not supported" >&2
        exit 3
        ;;
 stop)
        # No-op
        ;;
  *)
        echo "Usage: odloadconf [start]" >&2
        exit 3
        ;;
esac

:
